# dotfiles/i3/.config/systemd/user/screen-blanking.service.d/screen-blanking.timer
[Unit]
Description=Screen blanking service timer for every minute
Requires=screen-blanking.service

[Timer]
OnBootSec=0min
OnCalendar=*:0/1

[Install]
WantedBy=default.target
